# packaging/common/openpaas/openpaas.prerm
#! /bin/sh -e

systemctl stop openpaas
systemctl disable openpaas
systemctl daemon-reload
